# layers-agl/meta-agl/meta-agl-profile-weston/recipes-config/agl-login-manager/files/user-config.service
[Unit]
Description=AGL user config

[Service]
ExecStart=/bin/ln -sf @DISPLAY_XDG_RUNTIME_DIR@/wayland-0 %t/

[Install]
WantedBy=default.target
