# layers-agl/meta-agl/meta-agl-profile-weston/recipes-config/agl-login-manager/files/user-config.path
[Unit]
Description=AGL user config unit path

[Path]
PathExists=@DISPLAY_XDG_RUNTIME_DIR@/wayland-0

[Install]
WantedBy=default.target
