# layers-agl/meta-agl/meta-agl-profile-weston/recipes-graphics/agl-desktop-config/files/user-weston-term.service
[Unit]
Description=Terminal for weston user
After=user-config.service
Requires=user-config.service

[Service]
Type=simple
ExecStart=/usr/bin/weston-terminal
Restart=always

[Install]
WantedBy=multi-user.target
